# playbooks/templates/kill-stuck-apt.timer.j2
[Unit]
Description=Run kill-stuck-apt script daily at 2:00 AM

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
